# nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn config.wsgi:application --bind 0.0.0.0:8000 --workers 3
directory=/opt/project
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO
user=root

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=AUTO
stderr_logfile=AUTO
user=root
